# rig_alert.crontab
*/15 * * * * $HOME/rig-alert/rig_alert.sh > $HOME/rig-alert/cron.log 2>&1
